# supervisord.conf
[supervisord]
nodaemon=true

[program:process-engine]
command=process-engine
